# meeting-extension/local-server/com.meetingextension.whisper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meetingextension.whisper</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_DIR__/bin/python</string>
        <string>__SCRIPT_PATH__</string>
        <string>--model</string>
        <string>small</string>
        <string>--compute</string>
        <string>int8</string>
        <string>--device</string>
        <string>cpu</string>
        <string>--language</string>
        <string>ja</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>9000</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SCRIPT_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/whisper.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/whisper.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
